# services/memcached/run
#!/bin/bash
exec /usr/bin/memcached -u nobody -v
