# services/php-fpm/run
#!/bin/bash
exec /usr/sbin/php5-fpm -F -y /config/nginx/nginx-fpm.conf
